# playbooks/files/update/crontab
*/15 * * * * /usr/local/sbin/update-packages mailto=lara@localhost log_file=/var/log/update-packages.log
*/15 * * * * /usr/sbin/freebsd-update cron
