# playbooks/files/crontab
*/15 * * * * /usr/bin/env pkg upgrade -qn > /var/log/pkg-upgrade
*/15 * * * * /usr/local/sbin/system-builder /usr/src
