# playbooks/files/crontab
*/15 * * * * /usr/bin/env pkg upgrade -n | mail -s "pkg upgrade" lara@localhost
*/15 * * * * freebsd-update cron
